# per-host-config/ubuntu-diffia/rc.local
#!/bin/bash

echo "Put local customizations here"
