# per-host-config/ubuntu/rc.local
#!/bin/bash

# Put local customizations here

echo 1 > /proc/sys/kernel/perf_event_paranoid
echo 0 > /proc/sys/kernel/kptr_restrict

exit 0
